# packaging/systemd/vigil-agent.service
[Unit]
Description=Seppia Vigil host metrics agent
Documentation=https://seppia.ai/docs/vigil/agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=vigil
Group=vigil

# `run` is an explicit alias of the bare daemon invocation; it makes
# the `ps` line and the journal greppable. NOTE: vigil-agent uses Go's
# `flag` package, which stops parsing flags at the first non-flag
# argument — so `run` MUST come AFTER `--config`, not before, or
# `--config` is parsed as a positional and the daemon refuses to start
# with `unexpected argument(s)`.
ExecStart=/usr/bin/vigil-agent --config /etc/seppia/vigil.yml run

# Lifecycle. Pairs with vigil-agent's own --drain-timeout (default 5s),
# documented in README.md → "Process lifecycle".
Restart=on-failure
RestartSec=5s
# Don't loop on operator errors. 1 = Config (bad token / URL / token revoked),
# 3 = Usage (bad CLI flag). 2 = Runtime (transient) is allowed to retry.
RestartPreventExitStatus=1 3
TimeoutStopSec=10s
KillSignal=SIGTERM
SendSIGKILL=yes

# Resource ceilings. Tiny on purpose — if vigil-agent ever needs more
# than this we have a bug, not a sizing problem.
MemoryMax=128M
TasksMax=64

# Hardening. Chosen to be safe across systemd >= 240 (RHEL 8, Debian 10+,
# Ubuntu 20.04+). More aggressive options (ProtectProc, SystemCallFilter,
# CapabilityBoundingSet) are deliberately not set here so the unit drops
# in cleanly on every supported distro; operators with a known target
# distro can layer them on via a drop-in at /etc/systemd/system/vigil-agent.service.d/.
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
SystemCallArchitectures=native

[Install]
WantedBy=multi-user.target
